# cron.d/asterisk_sounds
0 3     * * *   root    cd /opt/asterisk_sounds; git pull origin master
